# deployment/supervisor.conf
[program:app]
command = uwsgi --ini /code/deployment/uwsgi.ini

[program:nginx]
command = service nginx restart
